# extra/cmder/config/profile.d/usb-stick-shell.cmd
@ECHO OFF

ECHO USB Stick Shell is not yet supported in BAT/CMD.  Use PowerShell instead
ECHO.

PAUSE

EXIT 0
